# crontab
30 0 * * * /usr/local/bin/python /app/HW_Scrapper.py >> /var/log/hw_scrapping.log 2>&1
30 0 * * * /usr/local/bin/python /app/W2TJ_Scrapper.py >> /var/log/w2tj_scrapping.log 2>&1
